# infra/avahi-photoframe.service
[Unit]
Description=Avahi mDNS alias - photoframe.local
After=avahi-daemon.service
Requires=avahi-daemon.service

[Service]
Type=simple
ExecStart=/bin/bash /opt/apps/photoframe/scripts/avahi-alias.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
